# potterymanager.conf
[program:potterymanager]
command=gunicorn -b localhost:8000 -w 4 potterymanager:app
directory=/opt/potterymanager
user=root
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true